# server/systemd/aiir-state-backup.service
[Unit]
Description=AIIR State Backup
After=network.target

[Service]
Type=oneshot
User=root
Group=www-data
ExecStart=/bin/bash -lc 'mkdir -p /var/backups/aiir-state && ts=$(date -u +%%Y%%m%%d-%%H%%M%%S) && tar -czf /var/backups/aiir-state/state-${ts}.tar.gz /var/www/aiir/ai/state && ls -1t /var/backups/aiir-state/state-*.tar.gz | tail -n +15 | xargs -r rm -f'
